# pulse.service
[Unit]
Description=Pulse — Project Status Tracker
After=network.target

[Service]
Type=simple
User=clawd
WorkingDirectory=/home/clawd/projects/pulse
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=18800
Environment=PULSE_DATA_DIR=/home/clawd/projects/pulse/data

[Install]
WantedBy=multi-user.target
